# catalyst-international/Documentation/Implementation/_Complete/position-monitor.service
# ============================================================================
# Name of Application: Catalyst Trading System
# Name of file: position-monitor.service
# Version: 1.0.0
# Last Updated: 2026-01-16
# Purpose: Systemd service for persistent HKEX position monitoring
#
# INSTALLATION:
#   1. Copy to: /etc/systemd/system/position-monitor.service
#   2. Run: systemctl daemon-reload
#   3. Run: systemctl enable position-monitor
#   4. Run: systemctl start position-monitor
#
# MANAGEMENT:
#   systemctl start position-monitor    # Start the service
#   systemctl stop position-monitor     # Stop the service
#   systemctl restart position-monitor  # Restart the service
#   systemctl status position-monitor   # Check status
#   journalctl -u position-monitor -f   # View live logs
# ============================================================================

[Unit]
Description=HKEX Position Monitor Service
Documentation=https://github.com/your-repo/catalyst-trading-system
After=network.target network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Working directory
WorkingDirectory=/root/Catalyst-Trading-System-International/catalyst-international

# Load environment from file
EnvironmentFile=/root/Catalyst-Trading-System-International/catalyst-international/.env

# Python virtual environment
Environment="PATH=/root/Catalyst-Trading-System-International/catalyst-international/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Main execution
ExecStart=/root/Catalyst-Trading-System-International/catalyst-international/venv/bin/python3 position_monitor_service.py

# Restart policy - always restart on failure
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=position-monitor

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/Catalyst-Trading-System-International/catalyst-international/logs

# Resource limits - keep it lightweight
MemoryMax=256M
CPUQuota=25%

# Watchdog - restart if unresponsive
WatchdogSec=600

[Install]
WantedBy=multi-user.target
