# catalyst-neural/catalyst-pipeline.timer
[Unit]
Description=Weekly neural model training pipeline

[Timer]
OnCalendar=Sun 20:00
Persistent=true

[Install]
WantedBy=timers.target
